# contrib/docker/w5base-ol9-runtime/supervisord.d/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND

stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=5000

autostart=false
autorestart=true
killasgroup=true
stopasgroup=true


