# contrib/docker/w5base-ol9-runtime/supervisord.d/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
redirect_stderr=true
stdout_logfile=/dev/null
stderr_logfile=/dev/stderr
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true


