# contrib/docker/w5base-ol9-runtime/supervisord.d/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/httpd.log
stderr_logfile=/var/log/supervisor/httpd_error.log

